# docker/php7/supervisord/conf.d/php-fpm.ini
[program:php-fpm7]
command = /usr/local/sbin/php-fpm -F
autostart=true
autorestart=true
priority=5
stderr_logfile=/var/log/supervisor/php-fpm7_error.log
stdout_logfile=/var/log/supervisor/php-fpm7_out.log